# extra/battery-control.service

[Unit]
Description=Battery Charge Control Service
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/battery-control.sh
Restart=always
RestartSec=10
User=root
Group=root
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security settings
NoNewPrivileges=true
ProtectSystem=false
ProtectHome=true
ReadWritePaths=/var/log /sys/class/power_supply /proc/acpi /run/tlp /tmp

[Install]
WantedBy=multi-user.target
